# etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

iperf3 -s &
/root/script/delayedmonstart.sh &
/root/script/flashled.sh 1
exit 0
