# etc/hotplug.d/button/buttons
logger "the button was ${BUTTON} and the action was ${ACTION}"
